# docker/supervisor/config/supervisord.conf
[supervisord]
nodaemon=true

[program:test]
# command=/test.sh # for testing only
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true